# ops/systemd/goat-ai.service
[Unit]
Description=GOAT AI FastAPI server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/GOAT_AI
ExecStartPre=/usr/bin/mkdir -p %h/GOAT_AI/var/logs %h/GOAT_AI/var/data
Environment=GOAT_HOST=127.0.0.1
Environment=GOAT_SERVER_PORT=62606
ExecStart=%h/GOAT_AI/.venv/bin/python -m uvicorn server:create_app \
    --factory \
    --host 127.0.0.1 --port 62606
Restart=always
RestartSec=5
StandardOutput=append:%h/GOAT_AI/var/logs/fastapi.log
StandardError=append:%h/GOAT_AI/var/logs/fastapi.log
Environment=GOAT_RUNTIME_ROOT=%h/GOAT_AI/var
Environment=GOAT_LOG_DIR=%h/GOAT_AI/var/logs
Environment=GOAT_LOG_PATH=%h/GOAT_AI/var/chat_logs.db
Environment=GOAT_DATA_DIR=%h/GOAT_AI/var/data
EnvironmentFile=-%h/GOAT_AI/.env

[Install]
WantedBy=default.target
